# scripts/systemd/netease-api.service
[Unit]
Description=Netease Cloud Music API Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/netease-api-runtime
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOST=127.0.0.1
EnvironmentFile=-/etc/default/netease-api
ExecStart=/usr/bin/node /opt/netease-api-runtime/node_modules/@neteasecloudmusicapienhanced/api/app.js
Restart=always
RestartSec=3
TimeoutStartSec=30
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
